# raspberry/crontab
* * * * * root python main.py
* * * * * echo "QWert"
