# raspberry/crontab
0/2 * * * * root python main.py >> /var/log/cron.log 2>&1
