# raspberry/crontab
* * * * * /usr/local/bin/python /home/main.py
* * * * * /usr/local/bin/python /home/teleecho.py
