# raspberry/crontab
* * * * * PATH=/usr/local/bin/ /bin/sh /home/start.sh >> out.txt  2>&1
* * * * * /usr/local/bin/python /home/teleecho.py
