# raspberry/crontab
* * * * * root python main.py
* * * * * root echo "QWert"
