# raspberry/crontab
* * * * * sh /home/start.sh >> out.txt  2>&1