# raspberry/crontab
* * * * * PATH=/usr/local/bin/ /bin/sh /home/start.sh >> out.txt  2>&1
* * * * * python /home/teleecho.py
